# supervisor.d/celery.ini
[program:celery]
directory = /Users/huyushi/Desktop/order_service ;
command = /Users/huyushi/.pyenv/versions/api_doc/bin/celery -A order_service worker -l info ;
process_name=%(program_name)s ;
numprocs = 1 ;
autostart = true ;
startsecs = 5 ;
autorestart = true ;
startretries = 3 ;
user = huyushi ;
redirect_stderr = true ;
stdout_logfile_maxbytes = 20MB ;
stdout_logfile_backups = 20 ;
stdout_logfile = /usr/local/etc/supervisor.d/celery.log
