# supervisor.d/web.ini
[program:uwsgi]
directory = /Users/huyushi/Desktop/order_service ;
command = /Users/huyushi/.pyenv/versions/api_doc/bin/uwsgi --ini uwsgi/uwsgi.ini ;
process_name=%(program_name)s ;
numprocs = 1 ;
stopsignal=INT ;
autostart = true ;
startsecs = 5 ;
autorestart = true ;
startretries = 3 ;
user = huyushi ;
redirect_stderr = true ;
stdout_logfile_maxbytes = 20MB ;
stdout_logfile_backups = 20 ;
stdout_logfile = /usr/local/etc/supervisor.d/web.log
